# doc/conf/etc/sv/freeradius/run
#!/bin/sh

sv check dbus >/dev/null || exit 1
exec /usr/sbin/freeradius -f -l stdout 2>&1
