# doc/conf/etc/sv/asterisk/run
#!/bin/sh

exec /usr/sbin/asterisk -cvvvvv 2>&1
