# doc/conf/etc/sv/nodejs/run
#!/bin/sh

exec /opt/node/bin/node /root/dev/int/src/server 2>&1
