# doc/conf/etc/sv/bind9/run
#!/bin/sh

exec /usr/sbin/named -f -u bind
