# doc/conf/etc/sv/hostapd/run
#!/bin/sh

exec /usr/sbin/hostapd /etc/hostapd/hostapd.conf 2>&1
